# kraken-news-rss.timer
[Unit]
Description=Kraken RSS News Collector Timer

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
