# kraken-analyzer.timer
[Unit]
Description=Run Kraken Analyzer every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
