# kraken-news-rss.service
[Unit]
Description=Kraken RSS News Collector
After=network.target kraken-collector.service

[Service]
Type=oneshot
User=USER_PLACEHOLDER
WorkingDirectory=WORKDIR_PLACEHOLDER
ExecStart=/usr/bin/python3 WORKDIR_PLACEHOLDER/news_collector_rss.py

[Install]
WantedBy=multi-user.target
